# science/mysos/vagrant/upstart/mesos-slave.conf
description "mesos slave"
start on stopped rc RUNLEVEL=[2345]
respawn
post-stop exec sleep 5

env MY_HOST=192.168.33.7
env ZK_HOST=192.168.33.7

exec /usr/local/sbin/mesos-slave \
    --master=zk://$ZK_HOST:2181/mesos/master \
    --ip=$MY_HOST \
    --hostname=$MY_HOST \
    --resources="cpus:4;mem:1024;disk:20000" \
    --isolation="cgroups/cpu,cgroups/mem" \
    --cgroups_enable_cfs

console log
