# vagrant/upstart/mysos.conf
description "mysos"
start on stopped rc RUNLEVEL=[2345]
respawn
post-stop exec sleep 5

exec /home/vagrant/mysos/vagrant/bin/mysos_scheduler.sh

console log
